# deploy/ai-assistant.service
[Unit]
Description=AI Elder Assistant Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-elder-assistant/backend
Environment="PATH=/opt/ai-elder-assistant/venv/bin"
ExecStart=/opt/ai-elder-assistant/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
